# vscode/rootfs/etc/s6-overlay/s6-rc.d/nginx-vscode/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: VS Code (Microsoft)
# Runs nginx reverse proxy in front of VS Code serve-web.
#
# HA Supervisor strips the ingress URL prefix before forwarding to the add-on.
# VS Code registers its /_vscode-cli/ handlers (including the mint-key AES key
# endpoint) under the full ingress path (--server-base-path).  This script
# substitutes the current ingress entry path into the nginx config template so
# that nginx can restore the prefix before forwarding to VS Code on port 1338.
#
# The ingress token rotates on every add-on restart.  VS Code sets a fresh
# vscode-secret-key-path cookie on the very next response, so the browser
# self-heals on the first page load after a restart.
# ==============================================================================
bashio::log.info 'Starting nginx reverse proxy for VS Code...'

ingress_entry=$(bashio::addon.ingress_entry)
bashio::log.info "Ingress entry: ${ingress_entry}"

sed "s|@@INGRESS@@|${ingress_entry}|g" \
    /etc/nginx/nginx-vscode.conf.tmpl \
    > /tmp/nginx-vscode.conf

exec nginx -c /tmp/nginx-vscode.conf 2>&1
